# server/deploy/riskybuyer.service
[Unit]
Description=RiskyBuyer API (FastAPI/uvicorn)
After=network.target

[Service]
User=riskybuyer
Group=riskybuyer
WorkingDirectory=/opt/riskybuyer/server
Environment=RB_DB=/var/lib/riskybuyer/riskybuyer.db
ExecStart=/opt/riskybuyer/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080 --proxy-headers --forwarded-allow-ips=*
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
